# deploy/tracker.service
[Unit]
Description=tracker — coordination store for coding agents
Documentation=https://github.com/chicagobuss/tracker
# Needs the network (binds LAN + ZeroTier IPs) and the Postgres container.
Wants=network-online.target
After=network-online.target docker.service zerotier-one.service
# Don't get locked into 'failed' while the ZeroTier/LAN IPs are still coming up.
StartLimitIntervalSec=0

[Service]
Type=simple
User=importer
Group=importer
WorkingDirectory=/path/to/tracker
EnvironmentFile=/path/to/tracker/.env
ExecStart=/path/to/tracker/tracker
# Retry until Postgres is reachable and the bound interfaces exist.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
